# kafka.conf
[supervisord]
nodaemon=true
user=root

[program:kafka]
command=/home/kafka/kafka_2.13-3.0.0/bin/start-kafka.sh
autostart=true
autorestart=true
user=root
