# kafka.conf
[program:kafka]
command=/home/kafka/kafka_2.13-2.6.0/bin/start-kafka.sh
autostart=true
autorestart=true
