# zookeeper.conf
[program:zookeeper]
command=/home/kafka/kafka_2.13-2.6.0/bin/start-zookeeper.sh
autostart=true
autorestart=true
